# root/etc/s6-overlay/s6-rc.d/svc-caddy/run
#!/usr/bin/with-contenv sh
exec caddy run --config /config/caddy/Caddyfile
